# scripts/autoguider1.root.crontab
# minute hour day_of_month month day_of_week command
# reboot machine 10:35 Mon-Fri
35 10 * * 1-5 /sbin/reboot
